# example/server_monitoring/service/rethinkdb/run
#!/bin/sh

exec rethinkdb -d /data/rethinkdb >> /var/log/monitor/rethinkdb.log 2>&1
