# example/server_monitoring/service/monitor/run
#!/bin/sh

cd /app && exec node index.js >> /var/log/monitor/node.log 2>&1
